# docker-compose.service
[Unit]
Description=%I – Docker Compose stack
PartOf=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%I                 # ← capital I = unescaped path
ExecStart=/usr/bin/docker compose --project-directory %I up -d --remove-orphans
ExecStop=/usr/bin/docker compose --project-directory %I down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target