# docker-compose.service
[Unit]
Description=%i – Docker Compose stack
PartOf=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%i
ExecStart=/usr/bin/docker compose --project-directory %i up -d --remove-orphans
ExecStop=/usr/bin/docker compose --project-directory %i down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target