# pete-panel.service
[Unit]
Description=Pete Panel - Docker Compose stack (16GB RAM / 4 vCPU)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/pete-panel

# 1) Always pull latest images for this stack
ExecStart=/usr/bin/docker compose -f docker-compose.prod-16ram-4cpu.yml pull

# 2) Rebuild images (in case Dockerfiles changed)
ExecStart=/usr/bin/docker compose -f docker-compose.prod-16ram-4cpu.yml build

# 3) Bring the stack up in detached mode
ExecStart=/usr/bin/docker compose -f docker-compose.prod-16ram-4cpu.yml up -d

# 4) Properly tear the stack down when stopping the service
ExecStop=/usr/bin/docker compose -f docker-compose.prod-16ram-4cpu.yml down

RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
