# wordpress-pete.service
[Unit]
Description=WordPress Pete - Docker Compose stack
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Absolute path to the folder that holds docker-compose.yml
WorkingDirectory=/opt/wordpress-pete
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes
TimeoutStartSec=0
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target